# dev/vm-builder/packer-local/ansible/roles/lvm-migration/files/lvm-migrate.service
[Unit]
Description=Migrate root filesystem to LVM on first boot
After=local-fs.target network.target
ConditionPathExists=/root/lvm-migrate.sh
ConditionPathExists=!/var/lib/lvm-migrate-done

[Service]
Type=oneshot
Environment="LVM_LV_ROOT_SIZE=8G"
Environment="DATA_LV_SIZE=2G"
ExecStart=/bin/bash /root/lvm-migrate.sh
ExecStartPost=/bin/touch /var/lib/lvm-migrate-done
ExecStartPost=/bin/systemctl disable lvm-migrate.service
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
